# examples/systemd/quantumLT-init.service
[Unit]
Description=QuantumLT initialize service

[Service]
Type=oneshot
ExecStart=/opt/quantumLT-cli/quantumLT-cli init -s %I